# deploy/d-brain-bot.service
[Unit]
Description=d-brain Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/agent-second-brain
ExecStart=/root/agent-second-brain/.venv/bin/python -m d_brain
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
